# dot_config/systemd/user/litellm-proxy.service
[Unit]
Description=LiteLLM Proxy - unified LLM API gateway
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=%h/.local/bin/litellm start --service
EnvironmentFile=-%h/.config/litellm/env
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
